# skel/.config/autostart/setxkbmap.desktop
[Desktop Entry]
Type=Application
Name=setxkbmap
Exec=setxkbmap pl
Comment=Set Polish keymap
